# @fsystem/src/etc.skel/.config/autostart/RSp - RspIdleManager.desktop
[Desktop Entry]
Exec=/usr/lib/rsp/proc/idleManagerStart
Hidden=true
Name=RSp - IdleManager
Type=Application
Version=1.0
