# @fsystem/src/etc.skel/.config/autostart/rsp-config-monitor-autostart.desktop
[Desktop Entry]
Comment=Autostart monitor settings for RSp-config-monitor
Exec=rsp-config-monitor -l
Name=rsp-config-monitor-autostart
OnlyShowIn=RSp
Type=Application
Version=1.0
